# telegram-bot.service
[Unit]
Description=Telegram Channel Discussion Defender
After=network.target

[Service]
Type=simple
User=YOUR_USERNAME
WorkingDirectory=/path/to/telegram-channel-discussion-defender
ExecStart=/path/to/telegram-channel-discussion-defender/telegram-channel-discussion-defender
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
